# packages/centos/tools-release.service
[Unit]
Description=tools release service

[Service]
Type=simple
ExecStart=/go/bin/tools-release

[Install]
WantedBy=multi-user.target
